# config/arbitrage_bots.target
[Unit]
Description=Arbitrage Bots
After=network.target

[Install]
WantedBy=multi-user.target
